# explorer.ini
[program:example]
command=/home/jovyan/explorer
directory=/home/jovyan
user=jovyan
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/jovyan/logs/explorer.std.log
loglevel=info

